# robot/system/robot.service
[Unit]
Description=Robot controller app

[Service]
Type=simple
ExecStart=/usr/bin/robot

[Install]
WantedBy=multi-user.target
Alias=robot.service

